# personal/roles/promtail/templates/promtail.service.j2
[Unit]
Description=Promtail log scraper
Requires=network-online.target

[Service]
ExecStart=/usr/local/bin/promtail -config.file=/etc/promtail/config.yml

[Install]
WantedBy=multi-user.target
